# APIPHP/run/systemd/system/systemd-sysctl.service.d/zzz-lxc-override.conf
[Service]
ExecStart=
ExecStart=-/lib/systemd/systemd-sysctl
