# APIPHP/run/systemd/system/rsync.service.d/zzz-lxc-service.conf
[Service]
ProtectControlGroups=no
ProtectKernelTunables=no
NoNewPrivileges=no
ProtectHome=no
ProtectSystem=no
PrivateDevices=no
PrivateTmp=no
ProtectKernelModules=no
ReadWritePaths=
